# debian/postrm
#!/bin/sh

#DEBHELPER#

set -e

gunicorn-debian stop bananas || true
rm -f /etc/gunicorn.d/bananas

rm -f /etc/nginx/sites-enabled/bananas
invoke-rc.d nginx reload || true
